# etc/systemd/system/greetd.service.d/render-fix.conf
[Service]
ExecStartPre=+/bin/sh -c "chgrp render /dev/dri/renderD* 2>/dev/null; chmod 0666 /dev/dri/renderD* 2>/dev/null; true"
