# dotfiles/config/systemd/user/matrix-wachter.service
[Unit]
Description=Matrix-Wächter — Update-Resets erkennen und selbstheilen

[Service]
Type=oneshot
ExecStart=%h/.local/bin/matrix-wachter
