# dotfiles/config/systemd/user/matrix-gesundheit.service
[Unit]
Description=Matrix Gesundheitsbericht
[Service]
Type=oneshot
ExecStart=%h/.local/bin/matrix-gesundheit
