# dotfiles/config/systemd/user/sync-kdeglobals.service
[Unit]
Description=kdeglobals aus DankMatugen-Schema aktualisieren
[Service]
Type=oneshot
ExecStart=%h/.local/bin/sync-kdeglobals.sh
