# dotfiles/config/systemd/user/theme-sonne.service
[Unit]
Description=Sonnenstands-Rhythmus planen (Hell/Dunkel nach echtem Tag)
[Service]
Type=oneshot
ExecStart=%h/.local/bin/matrix-sonne
