# dotfiles/config/systemd/user/theme-sonne-login.service
[Unit]
Description=Sonnenplanung bei Anmeldung (Catch-up)
[Service]
Type=oneshot
ExecStart=%h/.local/bin/matrix-sonne
[Install]
WantedBy=default.target
