# dotfiles/config/systemd/user/theme-light.service
[Unit]
Description=DMS auf light-Theme schalten
[Service]
Type=oneshot
ExecStart=/usr/bin/dms ipc call theme light
