# dotfiles/config/systemd/user/theme-dark.service
[Unit]
Description=DMS auf dark-Theme schalten
[Service]
Type=oneshot
ExecStart=/usr/bin/dms ipc call theme dark
