# etc/systemd/system/matrix-home-snapshot.timer
[Unit]
Description=Täglicher Matrix Home-Snapshot

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
