# etc/systemd/system/matrix-home-snapshot.service
[Unit]
Description=Matrix Home-Snapshot (btrfs Zeitmaschinen-Leitbild)
ConditionPathIsMountPoint=/home

[Service]
Type=oneshot
ExecStart=/usr/local/bin/matrix-home-snapshot
Nice=15
IOSchedulingClass=idle
